# systemd/azazel-epd.service.d/10-portal-detect.conf
[Unit]
Wants=network-online.target
After=network-online.target dhcpcd.service

[Service]
ExecStartPost=/bin/bash /usr/local/bin/portal_detect.sh
