# systemd/azazel-portal-viewer.service
[Unit]
Description=Azazel Captive Portal Viewer (Chromium + noVNC)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/usr/local/bin/azazel-portal-viewer.sh

[Service]
Type=simple
User=azazel
Group=azazel
WorkingDirectory=/home/azazel
EnvironmentFile=-/etc/default/azazel-gadget
EnvironmentFile=-/etc/default/azazel-zero
EnvironmentFile=-/etc/azazel-gadget/portal-viewer.env
EnvironmentFile=-/etc/azazel-zero/portal-viewer.env
Environment="PORTAL_RUNTIME_DIR=/run/azazel-portal-viewer"
RuntimeDirectory=azazel-portal-viewer
RuntimeDirectoryMode=0750
ExecStart=/usr/local/bin/azazel-portal-viewer.sh
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=azazel-portal-viewer
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
