# systemd/azazel-control-daemon.service
[Unit]
Description=Azazel-Zero Control Daemon (Unix socket action executor)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/home/azazel/Azazel-Zero/py/azazel_control/daemon.py

[Service]
Type=simple
User=root
WorkingDirectory=/home/azazel/Azazel-Zero
ExecStart=/usr/bin/python3 /home/azazel/Azazel-Zero/py/azazel_control/daemon.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=azazel-daemon

[Install]
WantedBy=multi-user.target
