# systemd/opencanary.service
[Unit]
Description=OpenCanary Honeypot (venv)
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
User=root
EnvironmentFile=-/etc/default/azazel-gadget
EnvironmentFile=-/etc/default/azazel-zero
Environment=VIRTUAL_ENV=/home/azazel/canary-venv
Environment=PATH=/home/azazel/canary-venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/local/bin/opencanary-start
ExecStop=/home/azazel/canary-venv/bin/opencanaryd --stop
PIDFile=/home/azazel/canary-venv/bin/opencanaryd.pid
Restart=always
RestartSec=3
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
